# codeready-workspaces-plugin-intellij/config/etc/supervisord.d/novnc.conf
[program:novnc]
command=websockify --web=/usr/share/novnc 8080 localhost:5900
stdout_logfile=/tmp/novnc.log
stderr_logfile=/tmp/novnc.log
user=user
autorestart=true
priority=500
