# codeready-workspaces-plugin-intellij/config/etc/supervisord.d/prevent-idle-timeout.conf
[program:prevent-idle-timeout]
command=/opt/prevent-idle-timeout.sh
stdout_logfile=/tmp/prevent-idle-timeout.log
stderr_logfile=/tmp/prevent-idle-timeout.log
user=user
autorestart=true
priority=500
