# codeready-workspaces-plugin-intellij/config/etc/supervisord.d/run-ide.conf
[program:run-ide]
environment=DISPLAY=":0"
command=/opt/run-ide.sh
stdout_logfile=/tmp/run-ide.log
stderr_logfile=/tmp/run-ide.log
user=user
autorestart=true
priority=400
