# codeready-workspaces-plugin-intellij/config/etc/supervisord.d/vncserver.conf
[program:vncserver]
environment=DISPLAY=":0"
command=/usr/libexec/vncserver :0
stdout_logfile=/tmp/vncserver.log
stderr_logfile=/tmp/vncserver.log
user=user
autorestart=false
priority=300
